# systemd/laju-go.service
[Unit]
Description=Renjana - Laju Go Production
After=network.target

[Service]
Type=simple
User=maulana
Group=maulana
WorkingDirectory=/home/maulana/projects/renjana
ExecStart=/home/maulana/projects/renjana/laju-go
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=renjana

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Environment
EnvironmentFile=/home/maulana/projects/renjana/.env

[Install]
WantedBy=multi-user.target
